# raspi-vdrtft/pi/systemd/user/browser.service

[Unit]

Description="Starting chromium on startup"
PartOf=default.target

[Service]

Environment="DISPLAY=:0"
ExecStart=/home/pi/bin/browser.sh

[Install]

WantedBy=default.target
