# raspi-vdrtft/pi/systemd/user/vdrosd.service

[Unit]

Description="Starting chromium on startup"
PartOf=default.target

[Service]

Environment="DISPLAY=:0"
ExecStart=/home/pi/bin/browser.sh vdrosd

[Install]

WantedBy=default.target
